# mining/vps/wechat-relay.service
[Unit]
Description=VoiceDrop WeChat publish relay (synchronous draft push from the whitelisted Tokyo IP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/wechat-relay
# WECHAT_RELAY_SECRET lives here (chmod 600, not in git). Must match the
# Cloudflare Pages secret of the same name.
EnvironmentFile=/opt/wechat-relay/relay.env
Environment=PORT=8848
# relay_server.py is stdlib-only and self-contained (no FILES_TOKEN, no R2). It runs
# on any Python 3.10+; the box already has 3.12 from the earlier deadsnakes install.
ExecStart=/usr/bin/python3.12 /opt/wechat-relay/relay_server.py
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/wechat-relay

[Install]
WantedBy=multi-user.target
